# scripts/com.fer.obsidian-rag-drift-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.fer.obsidian-rag-drift-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fer/repositories/obsidian-rag/.venv/bin/python</string>
        <string>/Users/fer/repositories/obsidian-rag/scripts/drift_watcher.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key><string>/Users/fer</string>
        <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/fer/.local/bin</string>
        <key>RAG_HOME</key><string>/Users/fer/.local/share/obsidian-rag</string>
        <key>NO_COLOR</key><string>1</string>
        <key>TERM</key><string>dumb</string>
    </dict>
    <key>StartInterval</key><integer>21600</integer>
    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><false/>
    <key>StandardOutPath</key><string>/Users/fer/.local/share/obsidian-rag/drift-watcher.log</string>
    <key>StandardErrorPath</key><string>/Users/fer/.local/share/obsidian-rag/drift-watcher.error.log</string>
</dict>
</plist>
